# worker.sample.ini
[program:ir-worker]
process_name=%(program_name)s_%(process_num)02d
command=cd /Volumes/Data/Sites/identifyrelevance/dashboard && /Users/uncovered/.mampenv/bin/php artisan queue:work --timeout=30
autostart=true
autorestart=true
user=uncovered
numprocs=32
redirect_stderr=true
stdout_logfile=/Volumes/Data/Sites/identifyrelevance/dashboard/storage/logs/%(program_name)s_%(process_num)02d.log
